# _debian/greenbone-security-assistant.after-install
#!/bin/bash
update-rc.d greenbone-security-assistant defaults
